# hh.conf
[program:hh]
command=python3 /home/curator2/hh/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/hh/test.err.log
stdout_logfile=/home/curator2/hh/test.out.log